# infra-repo-master/infra-repo-master/roles/mount/templates/systemctl.mount.j2
[Unit]
Description={{ item.description }}

[Mount]
What={{ item.partition }}
Where={{ item.mount_point }}
Type={{ item.fs_type }}
Options={{ item.options }}

[Install]
WantedBy=multi-user.target
